# kiwi/root/etc/profile.d/mono-mwf.sh
if [ ! -e /usr/share/doc/corefonts/EULA.html ]; then
  export MONO_MWF_SCALING=disable
fi
