# kiwi/root/etc/profile.d/mono-mwf.csh
if ( ! -e /usr/share/doc/corefonts/EULA.html ) then
  setenv MONO_MWF_SCALING disable
endif
